# src/main/systemd/zanzibar-backup.timer
[Unit]
Description=Zanzibar Backup Timer

[Timer]
Unit=zanzibar-backup.service
OnCalendar=*-*-* 01:45:00

[Install]
WantedBy=timers.target
